# systemd/sprouttrack-mqtt-exporter.service.template
[Unit]
Description=Sprout Track -> Home Assistant MQTT Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sprouttrack-mqtt-exporter
EnvironmentFile=/opt/sprouttrack-mqtt-exporter/.secrets.env
ExecStart=/opt/sprouttrack-mqtt-exporter/.venv/bin/python -m sprouttrack_exporter --config /opt/sprouttrack-mqtt-exporter/config.yaml --secrets /opt/sprouttrack-mqtt-exporter/.secrets.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
